# config/crontab
0 0 * * * cd /var/www/Challenge-S1 && php bin/console app:add-quantity
0 0 * * * cd /var/www/Challenge-S1 && php bin/console app:reminder-mail
